# deploy/home-computer/nightly-system-reboot.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Reboot the computer every day at 9:00 PM local time.
0 21 * * * root /sbin/shutdown -r now
